# fabrixe/deploy/fabrixe.service
[Unit]
Description=Fabrixe Enterprise Infrastructure Management Platform
Documentation=https://github.com/fabrixe/fabrixe
After=network-online.target avahi-daemon.service
Wants=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=fabrixe
Group=fabrixe
ExecStart=/opt/fabrixe/bin/fabrixe -config /etc/fabrixe/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fabrixe

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/fabrixe /etc/fabrixe
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
